# config/ultravox-elan.service
# Example systemd service for the ultravox-elan detector.
#
# Install:
#   sudo cp ultravox-elan.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ultravox-elan
#
# Logs:
#   journalctl -u ultravox-elan -f

[Unit]
Description=UltraVox ELAN live USV detector
After=network.target sound.target

[Service]
ExecStart=/opt/detector/ultravox-elan /opt/detector/ELAN.UVL --log-target 127.0.0.1:9999
WorkingDirectory=/opt/detector
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
